# infra/cron/refresh-sample-hash.timer
[Unit]
Description=Hourly refresh of /sample-hash.txt

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h
Persistent=true
Unit=refresh-sample-hash.service

[Install]
WantedBy=timers.target
